# update.service
[Unit]
Description=FIO Update Service

[Service]
ExecStart=/sysroot/home/fio/update.sh

[Install]
WantedBy=multi-user.target
